# .config/systemd/user/backup.service
[Unit]
Description=Run a backup
Requires=mnt-borg2.mount

[Service]
Type=oneshot
ExecStart=/home/nick/bin/backup
