# .config/systemd/user/backup.check.service
[Unit]
Description=Run a backup prune and check
Requires=mnt-borg2.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borgmatic -v2 prune check
