# .config/systemd/user/backup.service
[Unit]
Description=Run a backup

[Service]
Type=oneshot
ExecStart=/home/nick/bin/backup
