# .config/systemd/user/backup.service
[Unit]
Description=Run a backup
Requires=mnt-borg.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/borgmatic create -v2
