# .config/systemd/user/backup.check.service
[Unit]
Description=Run a backup check

[Service]
Type=oneshot
ExecStart=/home/nick/bin/backup.check
