# .config/systemd/user/backup.timer
[Unit]
Description=Run a backup

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
