# config/systemd/user/sway-session.target
[Unit]
Description=Current sway session
Requires=basic.target
BindsTo=graphical-session.target
